# Library/LaunchAgents/aria2c.rpc.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>aria2c.rpc</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/aria2c</string>
        <string>--enable-rpc</string>
        <string>-s5</string>
        <string>-x5</string>
        <string>-c</string>
        <string>-V</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/user/Downloads/aria2c</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
